# pi/systemd/greenhouse-ap.service
[Unit]
Description=Greenhouse AP mode (first-time WiFi setup)
After=NetworkManager.service greenhouse-firstboot.service
Wants=NetworkManager.service
ConditionPathExists=!/etc/greenhouse/.wifi_configured

[Service]
Type=oneshot
ExecStart=/home/pi/greenhouse/scripts/ap_up.sh
RemainAfterExit=yes
StandardOutput=journal

[Install]
WantedBy=multi-user.target
