# pi/systemd/greenhouse-firstboot.service
[Unit]
Description=Greenhouse first-boot provisioning
After=network-pre.target

[Service]
Type=oneshot
ExecStart=/home/pi/greenhouse/scripts/first_boot.sh
RemainAfterExit=yes
StandardOutput=journal

[Install]
WantedBy=sysinit.target
