# pi/systemd/greenhouse-firstboot.service
[Unit]
Description=Greenhouse first-boot provisioning
DefaultDependencies=no
Before=sysinit.target
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/home/pi/greenhouse/scripts/first_boot.sh
RemainAfterExit=yes
StandardOutput=journal

[Install]
WantedBy=sysinit.target
