# pi/systemd/greenhouse-portal.service
[Unit]
Description=Greenhouse pairing portal
After=network.target greenhouse-firstboot.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/greenhouse/portal
ExecStart=/usr/bin/python3 /home/pi/greenhouse/portal/portal.py
ExecStartPost=/sbin/iptables -t nat -A PREROUTING -p tcp --dport 8080 -j REDIRECT --to-port 80
Restart=always
RestartSec=5
StandardOutput=journal

[Install]
WantedBy=multi-user.target
